# debian/postinst
#systemctl enable com.mdm.app.event.service
chmod u+s /usr/sbin/MdmAppEvent
